# data/systemd/hawaii-notifications-daemon.service.in
[Unit]
Description=Hawaii notifications daemon
Requires=hawaii.service
After=hawaii.service

[Service]
Type=notify
NotifyAccess=main
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/hawaii-notifications-daemon

[Install]
WantedBy=hawaii.target
