# data/systemd/hawaii.target
[Unit]
Description=Hawaii session
#Requires=dbus.socket
AllowIsolate=true

[Install]
Alias=default.target
