# data/systemd/hawaii-kms.service.in
#
# This file is part of Hawaii.
#
# Copyright (C) 2014-2015 Pier Luigi Fiorini <pierluigi.fiorini@gmail.com>
#
# This program is free software; you can redistribute it and/or modify it
# under the terms of the GNU Lesser General Public License as published by
# the Free Software Foundation; either version 2.1 of the License, or
# (at your option) any later version.
#

[Unit]
Description=Hawaii (KMS)

[Service]
Type=notify
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus" QT_QPA_EGLFS_INTEGRATION=eglfs_kms
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/hawaii-session --mode=eglfs

[Install]
Alias=hawaii.service
WantedBy=hawaii.target
