# data/systemd/liri-shell-helper.service.in
[Unit]
Description=Liri Shell helper program
PartOf=liri-shell-helper.target
After=liri-shell.target
# Don't manually start/stop
RefuseManualStart=on
RefuseManualStop=on
# Force a session shutdown when this exits for a logout or a crash,
# because we are a display server
#OnFailure=liri-session-shutdown.service
#OnFailureJobMode=replace-irreversibly
#CollectMode=inactive-or-failed

[Service]
Type=notify
ExecStart=@ABSOLUTE_LIBEXECDIR@/liri-shell-helper
# Don't restart when aborting for a broken Wayland connection
Restart=on-watchdog
WatchdogSec=20
