# data/systemd/liri-shell-helper.service.in
[Unit]
Description=Liri Shell helper program
PartOf=liri-shell-helper.target
After=liri-shell.target
# Don't manually start/stop
RefuseManualStart=on
RefuseManualStop=on

[Service]
Type=notify
ExecStart=@ABSOLUTE_LIBEXECDIR@/liri-shell-helper
# Don't restart when aborting for a broken Wayland connection
Restart=on-abnormal
RestartPreventExitStatus=SIGABRT
WatchdogSec=20
